# qr3_29.conf
[program:qr3_29]
command=python3 /home/curator3/qr3_29/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_29/test.err.log
stdout_logfile=/home/curator3/qr3_29/test.err.log
